# deploy/systemd/signalops-marketops-daily.timer
[Unit]
Description=Run SignalOps MarketOps workflow after the U.S. market close

[Timer]
OnCalendar=Mon..Fri *-*-* 18:01:55 America/New_York
Persistent=true
AccuracySec=1s
Unit=signalops-marketops-daily.service

[Install]
WantedBy=timers.target
